# kiosk-auto-update.service
[Unit]
Description=Kiosk Auto-Update Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/kiosk-with-bob
ExecStart=/bin/bash /home/pi/kiosk-with-bob/auto-update.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target